# files/etc/uci-defaults/55_start_wifi_later
#!/bin/sh
/sbin/start_wifi_later &
exit 0
